# .config/systemd/user/face-security.service
[Unit]
Description=Face Detection Security Check
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/bin/face_check.sh
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
StandardOutput=journal
StandardError=journal
